# react-dev.service
[Unit]
Description=React Dev Application
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/react-dev
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/node /var/www/react-dev/server.js
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=react-dev

# Безопасность
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/react-dev

[Install]
WantedBy=multi-user.target

